# services/keylime_verifier.service.template
[Unit]
Description=The Keylime verifier service
After=network.target

[Service]
ExecStart=KEYLIMEDIR/keylime_verifier

[Install]
WantedBy=default.target
