# services/keylime_agent.service
[Unit]
Description=The Keylime compute agent

[Service]
ExecStart=/usr/bin/keylime_agent
KillSignal=SIGINT

[Install]
WantedBy=default.target
