# services/keylime_agent.service.template
[Unit]
Description=The Keylime compute agent

[Service]
ExecStart=KEYLIMEDIR/keylime_agent
KillSignal=SIGINT

[Install]
WantedBy=default.target
