# services/keylime_registrar.service.template
[Unit]
Description=The Keylime registrar service
After=network.target

[Service]
ExecStart=KEYLIMEDIR/keylime_registrar

[Install]
WantedBy=default.target
