# services/keylime_verifier.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/keylime_verifier

[Install]
WantedBy=default.target
