# ima_stub_service/tpm_emulator.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/tpm_with_ima.sh

[Install]
WantedBy=default.target